# gpu-fan-control.service
[Unit]
Description=GPU Fan Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gpu-fan-control.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target