# src/gpu-fan-control.service
[Unit]
Description=GPU to System Fan Control
After=multi-user.target

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/gpu-fan-control.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target